# ops-worker.service
[Unit]
Description=ops-worker machine state monitor
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/ops-worker -service -config /etc/ops-worker/config.yaml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
